# package/wwan/driver/quectel-mhi-pcie/files/pcie_mhi
#!/bin/sh /etc/rc.common

START=70

start() {
	/sbin/modprobe -q pcie_mhi
}

stop() {
	[ -d "/sys/module/pcie_mhi" ] && /sbin/rmmod pcie_mhi
}
