# package/wwan/app/luci-app-sms-tool/root/etc/uci-defaults/start-smsled
#!/bin/sh
# Copyright 2020 Rafał Wabik (IceG) - From eko.one.pl forum
# Licensed to the GNU General Public License v3.0.

/etc/init.d/smsled disable

exit 0
