# package/public/luci-app-athena-led/root/etc/uci-defaults/luci-app-athena-led
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@athena_led[-1]
	add ucitrack athena_led
	set ucitrack.@athena_led[-1].init=athena_led
	commit ucitrack
EOF

rm -rf /tmp/luci-*cache
exit 0
