# target/linux/bcm53xx/base-files/etc/init.d/vlan_gro
#!/bin/sh /etc/rc.common

START=90

boot() {
	stop
}

start() {
	ethtool -K eth0 gro on
}

stop() {
	ethtool -K eth0 gro off
}
